# apps/node-exporter/node-exporter.container
[Unit]
Description=node-exporter

[Container]
ContainerName=node-exporter
Environment=TZ=America/New_York
Exec=--path.rootfs=/host
Image=quay.io/prometheus/node-exporter:v1.8.2@sha256:4032c6d5bfd752342c3e631c2f1de93ba6b86c41db6b167b9a35372c139e7706
Network=host
Volume=/:/host:ro,rslave

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
